# etc/update-motd.d/00-header
#!/bin/sh

# Call the color file
. /etc/update-motd.d/colors

# Print the hostname with figlet
printf "\n"$GREEN
figlet $(hostname -s)
printf $NONE
